# provisioners/file/etc/systemd/system/jenkins-proxy.socket

# /etc/systemd/system/jenkins-proxy.socket

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
